# scripts/launchd/com.clipremote.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template only. Prefer:
    clipremote install-service
  which writes ~/Library/LaunchAgents/com.clipremote.daemon.plist with the real binary path,
  HOME, PATH, logs, RunAtLoad, and KeepAlive.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clipremote.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_CLIPREMOTE_PATH</string>
    <string>daemon</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
